# src/etc/cron.d/piguard
# Pi-guard: Safe network for home
# (c) 2021 Pi-guard (https://github.com/homesafeguard/pi-guard)

# Pi-guard: Night shutdown (energy saving)
55  23  *   *   *   root    PATH="$PATH:/usr/sbin:/usr/local/bin/" piguard shutdown > /var/log/piguard/debug/shutdown.log

# Pi-guard: Update to latest version
0  12   *   *   1   root    PATH="$PATH:/usr/sbin:/usr/local/bin/" piguard self-update > /var/log/piguard/debug/self-update.log
